# _config/systemd/user/obex.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/obexd -a
